# .config/systemd/user/backup-databases.service
[Unit]
Description=Database backup

[Service]
Type=oneshot
ExecStart=/disks/Files/Scripts/backup/backup-databases.sh

