# .config/systemd/user/backup-daily.service
[Unit]
Description=Daily backup

[Service]
Type=oneshot
ExecStart=/disks/Files/Scripts/backup/backup-daily.sh

