# .config/systemd/user/podcasts-sync.service
[Unit]
Description=Sync podcasts

[Service]
Type=oneshot
ExecStart=/disks/Files/Scripts/podcast/sync.sh

